# deploy/systemd/mission-control-discord-thread-ingest.timer
[Unit]
Description=Poll Discord thread for bounded Mission Control ingest payloads every two minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=120s
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
